# Hybrid_GPL_63168/DT-W724V-20140311/output/target/etc/init.d/fonap
#! /bin/sh
### BEGIN INIT INFO
# Provides:          skeleton
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Example initscript
# Description:       This file should be used to construct scripts to be
#                    placed in /etc/init.d.
### END INIT INFO

# Author: Foo Bar <foobar@baz.org>
#
# Please remove the "Author" lines above and replace them
# with your own name if you copy and modify this script.

# Do NOT "set -e"

# PATH should only include /usr/* if it runs after the mountnfs.sh script
PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="FON captive portal, dns and dhcp server"
NAME=fonap
DAEMON=/usr/sbin/$NAME
DAEMON_ARGS="--options args"
PIDFILE=/var/run/$NAME.pid
SCRIPTNAME=/etc/init.d/$NAME
# Load the VERBOSE setting and other rcS variables
#. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.0-6) to ensure that this file is present.
#. /lib/lsb/init-functions

if [ -f /etc/default/fonap ] ; then
   . /etc/default/fonap
fi

mkdir -p /var/fon

startfonap(){
    if [ "$START_FONAP" != "1" ] ; then
	echo "Fongw default off. Look at /etc/default/fonap"
	exit 0
    fi

	/etc/init.d/fonsmcd start &
}


stopfonap(){
	/etc/init.d/chilli stop &
	/etc/init.d/fonsmcd stop &
	foncontrl wifidown
	echo "[STOPFONAP]FON SSID DOWN"
}

stopfonservice(){
	killall pppd
	sleep 2
	killall chilloutd
	killall -SIGINT xl2tpd
	killall chillout-client	
	stopfonap
	killall fonsmcd
	kill $(ps aux|grep chilli|grep -v grep|grep -v radconfig|grep -v cron.daily|awk '{print $1}')
}


case "$1" in
  start)
	stopfonservice
	sleep 2
	/usr/sbin/chilloutd -r /usr/sbin/rule.sh -n &
	if [ $FONAP_TUNNELS = 1 ]; then
		mkdir -p /var/run/xl2tpd
		sleep 1
		/usr/sbin/xl2tpd -c /etc/xl2tpd.conf -D&
	fi
	sleep 1
	startfonap
	;;
  stop)
	stopfonservice
	;;
  #reload|force-reload)
	#
	# If do_reload() is not implemented then leave this commented out
	# and leave 'force-reload' as an alias for 'restart'.
	#
	#log_daemon_msg "Reloading $DESC" "$NAME"
	#do_reload
	#log_end_msg $?
	#;;
  restart|force-reload)
	#
	# If the "reload" option is implemented then remove the
	# 'force-reload' alias
	#
	stopfonap
	startfonap
	;;
  *)
	#echo "Usage: $SCRIPTNAME {start|stop|restart|reload|force-reload}" >&2
	echo "Usage: $SCRIPTNAME {start|stop|restart|force-reload}" >&2
	exit 3
	;;
esac

